# systemd.service
[Unit]
Description = login Service
[Service]
User=root
Environment=ZIPKIN_URL=http://172.31.27.189:9411/api/v2/spans
Environment=USERS_API_ADDRESS=172.31.18.169
ExecStart=/root/todo-auto/login/login
SyslogIdentifier=login

[Install]
WantedBy=multi-user.target
